# home-ether.mount
[Unit]

[Mount]
What=/dev/disk/by-id/google-ethereum-node-data
Where=/home/ether
Options=noatime

[Install]
WantedBy=multi-user.target
